# config/systemd/user/mpdscribble.service
[Unit]
Description = mpdscribble

[Service]
ExecStart = /usr/bin/mpdscribble -D
Restart = always

[Install]
WantedBy = default.target
