# config/systemd/user/mpdscribble.service
[Unit]
Description = mpdscribble
After=mpd.service

[Service]
ExecStart = /usr/bin/mpdscribble -D --daemon-user %u
Restart = always

[Install]
WantedBy = default.target
