# Library/LaunchAgents/eclimd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>eclimd</string>

        <key>ServiceDescription</key>
            <string>Eclim Daemon</string>
        <key>StandardErrorPath</key>
            <string>/Users/mfeltner/.eclimd.stderr</string>
        <key>StandardOutPath</key>
            <string>/Users/mfeltner/.eclimd.stdout</string>
        <key>UserName</key>
            <string>mfeltner</string>

        <key>OnDemand</key>
        <false/>

        <key>RunAtLoad</key>
        <true/>

        <key>EnvironmentVariables</key>
        <dict>
            <key>JAVA_HOME</key>
            <string>/Library/Java/Home</string>
        </dict>

        <key>ProgramArguments</key>
            <array>
                <string>/Users/mfeltner/Applications/eclipse/eclimd</string>
                <string>-XX:PermSize=2048</string>
                <string>-Xms1024M</string>
                <string>-Xmx2048M</string>
            </array>
    </dict>
</plist>
